# Tracker/tracker.service
[Unit]
After=mysql.service network-online.target

[Service]
ExecStart=/usr/local/bin/xbt_tracker

[Install]
WantedBy=default.target
