# Tracker/tracker.socket
[Socket]
DeferAcceptSec=30
ListenStream=2710

[Install]
WantedBy=sockets.target
